# deploy/command-center-relay.service
[Unit]
Description=Command Center Relay Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ccrelay
Group=ccrelay
WorkingDirectory=/var/lib/command-center-relay
ExecStart=/opt/command-center-relay/command-center-relay serve --port 9876 --host 127.0.0.1 --db /var/lib/command-center-relay/relay.db
Restart=on-failure
RestartSec=3
StartLimitIntervalSec=60
StartLimitBurst=5
Environment=RUST_LOG=info
LimitNOFILE=65536

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/command-center-relay
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
RestrictRealtime=true
MemoryMax=512M
TasksMax=256

[Install]
WantedBy=multi-user.target
